# deploy/systemd/zwmp-api.service
[Unit]
Description=ZWMP API
After=network.target

[Service]
Type=simple
User=zwmp
Group=zwmp
WorkingDirectory=/opt/zwmp
EnvironmentFile=/etc/zwmp/zwmp.env
ExecStart=/opt/zwmp/apps/api/.venv/bin/uvicorn zwmp_api.main:app --host 127.0.0.1 --port 8000 --proxy-headers
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/zwmp/data /opt/zwmp/.logs

[Install]
WantedBy=multi-user.target
